# tools/merovingian/daemon/monetdbd.service.in
[Unit]
Description=MonetDB database server daemon
Documentation=man:monetdbd https://www.monetdb.org/Documentation/monetdbd-man-page
After=network.target

[Service]
Type=forking
User=monetdb
Group=monetdb
ExecStart=@Ebindir@/monetdbd start @Elocalstatedir@/monetdb5/dbfarm
ExecStop=@Ebindir@/monetdbd stop @Elocalstatedir@/monetdb5/dbfarm
Restart=on-failure
PIDFile=@Erundir@/merovingian.pid
PrivateDevices=yes
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
